# mimir/rootfs/etc/s6-overlay/s6-rc.d/mimir/run
#!/command/with-contenv bash
# shellcheck shell=bash
# ==============================================================================
# Mímir Service
# ==============================================================================

# Configuration file path
CONFIG_PATH=/data/options.json

echo "=========================================="
echo "Mímir - Intelligent Home Assistant Agent"
echo "=========================================="

# Read configuration using jq
if [ -f "$CONFIG_PATH" ]; then
    export MIMIR_LLM_PROVIDER=$(jq -r '.llm_provider // "anthropic"' "$CONFIG_PATH")
    export MIMIR_LLM_API_KEY=$(jq -r '.llm_api_key // ""' "$CONFIG_PATH")
    export MIMIR_LLM_MODEL=$(jq -r '.llm_model // "claude-sonnet-4-20250514"' "$CONFIG_PATH")
    export MIMIR_LLM_BASE_URL=$(jq -r '.llm_base_url // ""' "$CONFIG_PATH")
    export MIMIR_TELEGRAM_OWNER_ID=$(jq -r '.telegram_owner_id // 0' "$CONFIG_PATH")
    export MIMIR_OPERATING_MODE=$(jq -r '.operating_mode // "normal"' "$CONFIG_PATH")
    export MIMIR_DEBUG=$(jq -r '.debug // false' "$CONFIG_PATH")
    export MIMIR_YOLO_MODE_DURATION=$(jq -r '.yolo_mode_duration_minutes // 10' "$CONFIG_PATH")
    export MIMIR_DELETIONS_PER_HOUR=$(jq -r '.deletions_per_hour // 5' "$CONFIG_PATH")
    export MIMIR_MODIFICATIONS_PER_HOUR=$(jq -r '.modifications_per_hour // 20' "$CONFIG_PATH")
    export MIMIR_GIT_ENABLED=$(jq -r '.git_enabled // true' "$CONFIG_PATH")
    export MIMIR_GIT_AUTHOR_NAME=$(jq -r '.git_author_name // "Mimir"' "$CONFIG_PATH")
    export MIMIR_GIT_AUTHOR_EMAIL=$(jq -r '.git_author_email // "mimir@asgard.local"' "$CONFIG_PATH")
else
    echo "Warning: Config file not found at $CONFIG_PATH"
fi

# Set log level based on debug setting
if [ "$MIMIR_DEBUG" = "true" ]; then
    export LOG_LEVEL="DEBUG"
    echo "Debug mode enabled"
else
    export LOG_LEVEL="INFO"
fi

echo "LLM Provider: ${MIMIR_LLM_PROVIDER}"
echo "LLM Model: ${MIMIR_LLM_MODEL}"
echo "Operating Mode: ${MIMIR_OPERATING_MODE}"
echo "Telegram Owner ID: ${MIMIR_TELEGRAM_OWNER_ID}"

# Debug: Show supervisor token status
if [ -n "$SUPERVISOR_TOKEN" ]; then
    echo "SUPERVISOR_TOKEN: [SET]"
else
    echo "SUPERVISOR_TOKEN: [NOT SET]"
fi

# Create data directories
mkdir -p /data/mimir

# Run the application
echo "Starting Python application..."
cd /opt/mimir
exec python3 -m app.main
